# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-ollama/run
#!/bin/sh
# Set environment variables for Ollama
export HOME=/root
export OLLAMA_MODELS=/data/ollama_models
export OLLAMA_HOST=0.0.0.0:${COCOSEARCH_OLLAMA_PORT:-11434}

echo "Starting Ollama..."

# Ensure model directory exists and has correct permissions
mkdir -p /data/ollama_models

# Auto-pull model if not present (supports override via COCOSEARCH_EMBED_MODEL)
MODEL="${COCOSEARCH_EMBED_MODEL:-nomic-embed-text}"

# Start Ollama with readiness notification
# -n 120: max 120 attempts, -s 2000: 2s between checks, -w 10000: 10s timeout per check
exec s6-notifyoncheck -d -n 120 -s 2000 -w 10000 \
    -c /etc/s6-overlay/s6-rc.d/svc-ollama/data/check \
    ollama serve
