# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-mcp/run
#!/bin/sh
echo "Starting CocoSearch MCP server..."

# Wait a moment for database connections to stabilize
sleep 1

# Resolve configuration from environment
# COCOSEARCH_MCP_TRANSPORT is our env var, MCP_TRANSPORT is what the CLI reads
TRANSPORT="${COCOSEARCH_MCP_TRANSPORT:-http}"
PORT="${COCOSEARCH_MCP_PORT:-3000}"

echo "Transport: $TRANSPORT, Port: $PORT"

# Start the MCP server
# Use exec to replace shell with python process for proper signal handling
exec python -m cocosearch mcp --transport "$TRANSPORT" --port "$PORT"
