# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-postgresql/run
#!/bin/sh
set -e

# Add PostgreSQL binaries to PATH
export PATH="/usr/lib/postgresql/17/bin:$PATH"

# Initialize PostgreSQL data directory if needed
if [ ! -f /data/pg_data/PG_VERSION ]; then
    echo "Initializing PostgreSQL data directory..."
    chown -R postgres:postgres /data/pg_data
    su postgres -c "initdb -D /data/pg_data"
fi

# Create cocosearch user and database if needed
if [ -f /data/pg_data/PG_VERSION ] && [ ! -f /data/pg_data/.cocosearch_initialized ]; then
    echo "Creating cocosearch database..."
    su postgres -c "pg_ctl start -D /data/pg_data -l /tmp/pg_init.log -w"
    su postgres -c "psql -c \"CREATE USER cocosearch WITH PASSWORD 'cocosearch';\"" || true
    su postgres -c "psql -c \"CREATE DATABASE cocosearch OWNER cocosearch;\"" || true
    su postgres -c "psql -d cocosearch -c \"CREATE EXTENSION IF NOT EXISTS vector;\"" || true
    su postgres -c "pg_ctl stop -D /data/pg_data -m fast"
    touch /data/pg_data/.cocosearch_initialized
fi

echo "Starting PostgreSQL..."
# Use s6-notifyoncheck for readiness notification
exec s6-notifyoncheck -d -n 60 -s 1000 -w 5000 \
    -c /etc/s6-overlay/s6-rc.d/svc-postgresql/data/check \
    su postgres -c "postgres -D /data/pg_data"
